# am4-fleet-node/systemd/am4-oxen-facade.service
[Unit]
Description=AM4 oxen OpenAI-compatible facade
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=derek
Group=derek
WorkingDirectory=/home/derek/am4-fleet-node
EnvironmentFile=-/home/derek/.config/am4-fleet/oxen.env
ExecStart=/usr/bin/python3 /home/derek/am4-fleet-node/scripts/oxen-facade.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

